# .devcontainer/services/nvim.service
[Unit]
Description=Headless Neovim server
After=network.target

[Service]
Type=simple
# Where the RPC socket will live (change if you prefer)
Environment=NVIM_LISTEN_ADDRESS=/run/user/%U/nvim/nvim.sock
# Clean up stale socket before starting
ExecStartPre=/usr/bin/rm -f /run/user/%U/nvim/nvim.sock
# Start headless server
ExecStart=/usr/bin/nvim --headless --listen ${NVIM_LISTEN_ADDRESS}
Restart=always
RestartSec=1
# Optional: set a working dir
WorkingDirectory=%h
# Optional: resource limits
Nice=5
OOMScoreAdjust=100

[Install]
WantedBy=default.target
